# variants/ap/etc/uci-defaults/99-ap-mode-network
#!/bin/sh
# uci-defaults: AP mode network configuration
# Self-deletes via rm (creates overlayfs whiteout).
# Whiteout path is added to sysupgrade.conf so the deletion
# state survives sysupgrade, preserving user modifications.

# Bridge all ports (lan1, lan2, lan3, wan)
uci -q delete network.device
uci set network.device='device'
uci set network.device.name='br-lan'
uci set network.device.type='bridge'
uci add_list network.device.ports='lan1'
uci add_list network.device.ports='lan2'
uci add_list network.device.ports='lan3'
uci add_list network.device.ports='wan'

# LAN: DHCP client
uci set network.lan='interface'
uci set network.lan.device='br-lan'
uci set network.lan.proto='dhcp'

# LAN static fallback
uci set network.lan_static='interface'
uci set network.lan_static.device='br-lan'
uci set network.lan_static.proto='static'
uci set network.lan_static.ipaddr='192.168.10.1'
uci set network.lan_static.netmask='255.255.255.0'

# IPv6: accept RA on br-lan (AP mode)
uci set network.lan.ipv6='auto'

uci commit network

# Self-delete and preserve whiteout across sysupgrade
rm "$0"
WHITOUT="/etc/uci-defaults/.wh.$(basename "$0")"
grep -qx "$WHITOUT" /etc/sysupgrade.conf 2>/dev/null || echo "$WHITOUT" >> /etc/sysupgrade.conf

exit 0
