# variants/ap/etc/uci-defaults/99-init-ap
#!/bin/sh
# uci-defaults: AP mode network configuration (first boot only)
# Skip if this looks like a sysupgrade (hostname changed or password set).
# OpenWrt automatically removes successfully executed uci-defaults scripts.

HOSTNAME=$(uci -q get system.@system[0].hostname)
if [ -n "$HOSTNAME" ] && [ "$HOSTNAME" != "OpenWrt" ]; then
    exit 0
fi
if grep -q '^root:[^!:]' /etc/shadow 2>/dev/null; then
    exit 0
fi

MAC_SUFFIX=$(tr -d ':' </sys/class/net/eth0/address | sed 's/.*\(....\)$/\1/' | tr 'a-f' 'A-F')

uci set system.@system[0].hostname="CMCC-$MAC_SUFFIX"

# Bridge all ports
uci -q delete network.@device[0]
uci set network.device='device'
uci set network.device.name='br-lan'
uci set network.device.type='bridge'
uci add_list network.device.ports='lan1'
uci add_list network.device.ports='lan2'
uci add_list network.device.ports='lan3'
uci add_list network.device.ports='wan'

# LAN: DHCP client
uci set network.lan='interface'
uci set network.lan.device='br-lan'
uci set network.lan.proto='dhcp'

# LAN static fallback
uci set network.lan_static='interface'
uci set network.lan_static.device='br-lan'
uci set network.lan_static.proto='static'
uci set network.lan_static.ipaddr='192.168.10.1'
uci set network.lan_static.netmask='255.255.255.0'

# IPv6: accept RA on br-lan (AP mode)
uci set network.lan.ipv6='auto'

# Remove wan/wan6 interfaces (all ports bridged in AP mode)
uci -q delete network.wan
uci -q delete network.wan6

uci commit network
uci commit system

exit 0
